# DIR875GPL/SRC/ubicom-distro-3.0/uClinux/vendors/Ubicom/DIR857RU/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS S start
::shutdown:/etc/init.d/rcS K stop
::sysinit:/bin/getty -nl /bin/sh 115200 $CONSOLE;

