# DIR657GPL/SRC/ubicom-distro/uClinux/vendors/Ubicom/DIR657G/rootfs/etc/init.d/dropcache
#!/bin/sh /etc/rc.common

START=490
STOP=10

start() {
	/etc/dropcache.sh &
}

stop() {
	killall dropcache.sh
}

