# DIR857GPL/SRC/ubicom-distro-3.0/uClinux/vendors/Ubicom/DIR857G/rootfs/etc/init.d/smp
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 Ubicom, Inc 

START=99

start() {
	/sbin/set_smp_affinity eth0 eth1 pciej pciek &
}

